# misc/after-staple.target
# -*- conf -*-

[Unit]
Description=After staple target
Requires=staple.service graphical.target
After=staple.service graphical.target
